# deployment_files/app.conf
[program:presence_dashboard]
directory=/var/flaskapp/simpleapp
command=/bin/bash -E -c ./start.sh
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
user=www-data